# picokey.service
[Unit]
Description=Picokey

[Service]
ExecStart=%h/.cargo/bin/picokey

[Install]
WantedBy=default.target
